# .scripts/caps2esc.service
[Unit]
Description=Remaps capslock to escape

[Service]
ExecStart=/home/nharsch/.scripts/caps2esc_start.sh

[Install]
WantedBy=multi-user.target
